# infra/backup/com.aihufu.backup.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aihufu.autosave-backup</string>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/REPO_PLACEHOLDER/infra/backup/autosave-and-backup.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/REPO_PLACEHOLDER</string>
  <key>StandardOutPath</key>
  <string>/REPO_PLACEHOLDER/infra/backup/.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/REPO_PLACEHOLDER/infra/backup/.launchd.err.log</string>
</dict>
</plist>
